# jobfilter.service
[Unit]
Description=Job Filter Pipeline Orchestrator
# Wait for network, Ollama, and a synced clock before starting. The clock
# matters: the Pi has no RTC battery, and jobs scored with a stale boot-time
# clock get timestamps behind the sync watermark — they'd be mistaken for
# hand-deleted rows and never reach the tracker.
After=network-online.target ollama.service time-sync.target
Wants=network-online.target ollama.service time-sync.target

[Service]
Type=simple
User=bluke
WorkingDirectory=/home/bluke/job_filter

ExecStart=/usr/bin/python3 /home/bluke/job_filter/scripts/orchestrator.py

# Restart on crash, OOM, or watchdog kill — but NOT on a clean systemctl stop.
# The orchestrator exits with code 0 on SIGTERM, so manual stops won't trigger a restart.
Restart=on-failure
RestartSec=15

# Caps CPU to 3 cores so the watchdog doesn't kill the process
CPUQuota=300%

StandardOutput=append:/home/bluke/job_filter/data/filter.log
StandardError=append:/home/bluke/job_filter/data/filter.log

[Install]
WantedBy=multi-user.target
