# jobfilter.service
[Unit]
Description=Job Filter Pipeline Orchestrator
# Wait for network and Ollama before starting
After=network-online.target ollama.service
Wants=network-online.target ollama.service

[Service]
Type=simple
User=bluke
WorkingDirectory=/home/bluke/job_filter

ExecStart=/usr/bin/python3 /home/bluke/job_filter/orchestrator.py

# Restart on crash, OOM, or watchdog kill — but NOT on a clean systemctl stop.
# The orchestrator exits with code 0 on SIGTERM, so manual stops won't trigger a restart.
Restart=on-failure
RestartSec=15

# Caps CPU to 3 cores so the watchdog doesn't kill the process
CPUQuota=300%

StandardOutput=append:/home/bluke/job_filter/filter.log
StandardError=append:/home/bluke/job_filter/filter.log

[Install]
WantedBy=multi-user.target
