# jobfilter.service
[Unit]
Description=Job Filter Run (one-shot until all jobs scored)
After=ollama.service
Wants=ollama.service

[Service]
Type=simple
User=bluke
WorkingDirectory=/home/bluke/job_filter
# Re-runs the filter; the seen-list makes each restart resume where it stopped.
# Exits 0 (success, no restart) only when the whole batch completes cleanly.
ExecStart=/usr/bin/python3 /home/bluke/job_filter/filter_jobs.py scraped_jobs.json --all
# If it dies for ANY reason (crash, Ollama 500 cascade, killed), wait and retry.
Restart=on-failure
RestartSec=15
StandardOutput=append:/home/bluke/job_filter/filter.log
StandardError=append:/home/bluke/job_filter/filter.log

[Install]
WantedBy=multi-user.target
